# server/init/weather-reader.conf
# weather-reader - weather-reader job file
# http://upstart.ubuntu.com/cookbook/

description "Weather Reader Daemon for reading 433MHz radio sensors"
author "Hasse Bjork<hassebjork@users.noreply.github.com>"

# Stanzas
#
# Stanzas control when and how a process is started and stopped
# See a list of stanzas here: http://upstart.ubuntu.com/wiki/Stanzas#respawn

# When to start the service
start on runlevel [2345]

# When to stop the service
stop on runlevel [016]

# Automatically restart process if crashed
respawn

# Essentially lets upstart know the process will detach itself to the background
expect fork

# Run before process
pre-start script
    [ -d /var/run/weather-reader ] || mkdir -p /var/run/weather-reader
    echo "Put bash code here"
end script

# Start the process
exec weather-reader